# kiosk/systemd/flipperklubben-kiosk.service
[Unit]
Description=Flipperklubben Kiosk (Chromium helskärm)
Documentation=https://github.com/flipperklubben/pin-kiosk
After=systemd-user-sessions.service network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kiosk
PAMName=login
TTYPath=/dev/tty1
# Ge sessionen en seat så att X kan startas utan root.
UnsetEnvironment=TERM
StandardInput=tty
StandardOutput=journal
StandardError=journal
TTYReset=yes
TTYVHangup=yes

# Starta X med kiosk-sessionen som enda klient.
ExecStart=/usr/bin/xinit /opt/flipperklubben-kiosk/bin/kiosk-session.sh -- :0 vt1 -keeptty -nolisten tcp

Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
